# app-admin/vaultwarden/files/init
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/var/lib/vaultwarden/vaultwarden"
pidfile="/run/vaultwarden.pid"
description="Vaultwarden - Unofficial Bitwarden compatible server written in Rust"
start_stop_daemon_args="--env VAULTWARDEN_CONFIG=$VAULTWARDEN_CONFIG --user ${VAULTWARDEN_USER} --group ${VAULTWARDEN_GROUP} --pidfile ${pidfile} --make-pidfile --background"

depend() {
	need net
	use mysql postgresql
}
