# app-admin/bitwarden_rs/files/init
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/var/lib/bitwarden_rs/bitwarden_rs"
pidfile="/run/bitwarden_rs.pid"
description="Unofficial Bitwarden compatible server written in Rust"
start_stop_daemon_args="--env BITWARDEN_RS_CONFIG=$BITWARDEN_RS_CONFIG --user ${BITWARDEN_RS_USER} --group ${BITWARDEN_RS_GROUP} --pidfile ${pidfile} --make-pidfile --background"

depend() {
	need net
	use mysql postgresql
}
